# container_files/etc/supervisor/conf.d/cron.conf
[program:cron]
command=/usr/sbin/cron -f
autorestart=true
; cron uses syslog
stdout_logfile=/dev/null
stderr_logfile=/dev/null
